# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# ------------------------
# Laravel Log Tail Process
# ------------------------
[program:laravel-log]
command=tail -f /app/storage/logs/laravel.log
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel-error.log
stdout_logfile=/var/log/laravel-output.log
user=www-data

# ------------------------
# PHP-FPM Process
# ------------------------
[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm-error.log
stdout_logfile=/var/log/php-fpm-output.log
user=www-data

# ------------------------
# Nginx Process
# ------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx-error.log
stdout_logfile=/var/log/nginx-output.log
user=root
